# server/systemd/disp-api.service
[Unit]
Description=Disp Registry FastAPI
After=network.target

[Service]
User=andrey
Group=andrey
WorkingDirectory=/home/andrey/opt/disp-app/server

Environment="DISP_DATA_DIR=/home/andrey/opt/disp-app/server/data"
Environment="DISP_ALLOWED_ORIGINS=https://botmaks.ru,https://www.botmaks.ru"

Environment="DISP_AUTH_USERNAME=admin"
Environment="DISP_AUTH_PASSWORD=897653897654Long"
Environment="DISP_SESSION_SECRET=7c4b8f91d2e64c0bb2d6b8f3f4fbb5b8a3a1c5d0e2f7a9c1d4e6f8b0c2a4d6e8"
Environment="DISP_SESSION_SECURE=true"
Environment="DISP_SESSION_TTL=28800"

ExecStart=/home/andrey/opt/disp-app/server/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
